# deploy/emlms.socket
[Unit]
Description=EMLMS Gunicorn socket

[Socket]
ListenStream=/run/emlms/emlms.sock
SocketUser=www-data

[Install]
WantedBy=sockets.target
